# systemd/opencode-trace-proxy.service
[Unit]
Description=OpenCode trace proxy
After=network-online.target

[Service]
WorkingDirectory=%h/.local/share/opencode-trace
ExecStart=/usr/bin/node %h/.local/share/opencode-trace/src/server.js
Environment=HOST=127.0.0.1
Environment=PORT=5010
Environment=UPSTREAM_BASE_URL=https://api.deepseek.com
Environment=RAW_TRACE_DIR=%h/.local/share/opencode-trace/logs
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
